# plugins/socket_backend/examples/launchd/io.github.telegramcodexbot.agent-node.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.telegramagentbot.agent-node</string>

  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USER/Projects/TelegramCodexBot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USER/Projects/TelegramCodexBot/.venv/bin/telegram-agent-node</string>
    <string>--node-id</string>
    <string>macbook</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8765</string>
    <string>--max-message-bytes</string>
    <string>26214400</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/Library/Logs/telegram-agent-node.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/Library/Logs/telegram-agent-node.err.log</string>
</dict>
</plist>
